# 00root/etc/init/nginx.conf
description "Nginx web server"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]

expect fork
respawn

exec /usr/local/nginx/sbin/nginx
