# 00root/etc/init/nginx.conf
description "Nginx web server"

start on (net-device-up and local-filesystems)
stop on runlevel [016]

expect fork
respawn
exec /usr/local/nginx/sbin/nginx
