# 00root/etc/init/vnstat.conf
description "vnStat network traffic monitor"

start on (net-device-up and local-filesystems)
stop on runlevel [016]

respawn
exec /var/server/vnstat/vnstatd -d
