# 00root/etc/init/vnstat.conf
description "vnStat network traffic monitor"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]

respawn

exec /var/server/vnstat/vnstatd -d
